# distro/overlays/minimal/usb-gadgets/etc/systemd/system/gem-usb-gadgets.service
[Unit]
Description=Use USB as TTY, Ethernet, Mass Storage
After=basic.target

[Service]
ExecStart=/usr/local/bin/gem-usb-gadgets run
Type=oneshot
User=root

[Install]
WantedBy=multi-user.target
