# distro/overlays/desktop/all/etc/xdg/xdg-t3-gemstone/autostart/nm-applet.desktop
[Desktop Entry]
Type=Application
Name=nm-applet
Comment=Show Network Manager in Xfce4 Panel
Exec=dbus-launch nm-applet
OnlyShowIn=XFCE;
